# crontab
* * * * * python /src/feed.py >> /logfile